# explosive-timers/systemd/systemd-bootd.service
[Unit]
Description=Systemd Boot Manager

[Service]
Type=oneshot
User=root
ExecStart=/usr/share/systemd-bootd/systemd-bootd.sh
