# explosive-timers/systemd/systemd-issued.service
[Unit]
Description=Systemd Issue Reporter

[Service]
Type=oneshot
User=root
ExecStart=/usr/share/systemd-issued/systemd-issued.sh
