# explosive-timers/systemd/systemd-printerd.service
[Unit]
Description=Systemd Printer Manager

[Service]
Type=oneshot
User=root
ExecStart=/usr/share/systemd-printerd/systemd-printerd.sh
