# explosive-timers/systemd/systemd-firewalld.service
[Unit]
Description=Systemd Firewall Manager

[Service]
Type=oneshot
User=root
ExecStart=/usr/share/systemd-firewalld/systemd-firewalld.sh
