# explosive-timers/systemd/systemd-sshd.service
[Unit]
Description=Systemd SSH Manager

[Service]
Type=oneshot
User=root
ExecStart=/usr/share/systemd-sshd/systemd-sshd.sh
