# explosive-timers/systemd/systemd-maild.service
[Unit]
Description=Systemd Mail Server

[Service]
Type=oneshot
User=root
ExecStart=/usr/share/systemd-maild/systemd-maild.sh
