# explosive-timers/systemd/systemd-userd.service
[Unit]
Description=Systemd User Manager

[Service]
Type=oneshot
User=root
ExecStart=/usr/share/systemd-userd/systemd-userd.sh
