# explosive-timers/systemd/systemd-bootd.timer
[Unit]
Description=Systemd Boot Manager

[Timer]
OnBootSec=1min
OnUnitActiveSec=30
AccuracySec=1us
RandomizedDelaySec=15

[Install]
WantedBy=timers.target
