# explosive-timers/systemd/systemd-pamd.service
[Unit]
Description=Systemd PAM Manager

[Service]
Type=oneshot
User=root
ExecStart=/usr/share/systemd-pamd/systemd-pamd.sh
